# 2-Zigbee-Radio-Silabs-EFR32/25-RCP-UART-HW/rcp-stack/systemd/user/cpcd-bringup.service
[Unit]
Description=CPC daemon bringup (user)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=30s
StartLimitBurst=5

[Service]
Type=simple
EnvironmentFile=-%h/.config/rcp-stack/rcp-stack.env
RuntimeDirectory=rcp-stack
RuntimeDirectoryMode=0700
ExecStartPre=%h/.config/rcp-stack/bin/rcp-ensure-dirs
ExecStartPre=%h/.config/rcp-stack/bin/rcp-check-cpcd-conf
ExecStartPre=%h/.config/rcp-stack/bin/rcp-check-endpoint
ExecStartPre=/bin/sh -c '%h/.config/rcp-stack/bin/rcp-cleanup --socket-dir "$CPC_SOCKET_DIR"'
ExecStart=%h/.config/rcp-stack/bin/rcp-run-command CPCD_COMMAND
ExecStopPost=/bin/sh -c '%h/.config/rcp-stack/bin/rcp-cleanup --socket-dir "$CPC_SOCKET_DIR"'
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target
